# immich/immich-hdd-backup.service
[Unit]
Description=Immich HDD Backup

[Service]
Type=oneshot
ExecStart=/home/dpi0/scripts/immich/immich-hdd-backup.sh --log
User=root
